# etc2/cron.weekly/backup
#!/bin/bash
# Run every other Sunday:
(( ($(date +%s) / 86400 - 3) % 14 )) || chronic /opt/jwodder/bin/drop-backup
